# deploy/systemd/railways-secretariat-backup.timer
# =============================================================================
# railways-secretariat-backup.timer
#
# Schedules railways-secretariat-backup.service every night at 02:30 local
# time, with a 30-minute randomised delay so concurrent backups across a
# fleet do not all hit the same disk at the same instant.
#
# Persistent=true ensures that if the host was off when the schedule fired,
# the missed run will execute on the next boot. Combined with the
# 14-archive retention this gives operators a generous, mostly hands-off
# safety net.
# =============================================================================
[Unit]
Description=Daily backup of Railway Secretariat System data
Documentation=file:///opt/railways-secretariat-flutter/docs/DISASTER_RECOVERY.md
Requires=railways-secretariat-backup.service

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=30min
Persistent=true
Unit=railways-secretariat-backup.service

[Install]
WantedBy=timers.target
