# deploy/systemd/railways-secretariat-backup.service
# =============================================================================
# railways-secretariat-backup.service
#
# Oneshot service that produces a consistent snapshot of the Railway
# Secretariat System (SQLite DB + uploaded attachments). Triggered by
# railways-secretariat-backup.timer; not enabled by default.
#
# Install:
#   sudo cp deploy/systemd/railways-secretariat-backup.{service,timer} \
#           /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now railways-secretariat-backup.timer
#
# The service runs as root because the backup script needs the docker
# socket. If you prefer to run it as a non-root user, add that user to the
# `docker` group and override `User=` here.
# =============================================================================
[Unit]
Description=Railway Secretariat System — consistent SQLite + attachments backup
Documentation=file:///opt/railways-secretariat-flutter/docs/DISASTER_RECOVERY.md
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Hardcoded to the standard production deployment path; override the path
# in /etc/systemd/system/railways-secretariat-backup.service.d/ if you
# install the repo somewhere else.
ExecStart=/usr/bin/env bash /opt/railways-secretariat-flutter/deploy/scripts/backup.sh

# Defaults; safe to override per-host via a drop-in.
Environment=VOLUME_NAME=railways_secretariat_data
Environment=BACKUP_DIR=/var/backups/railways-secretariat
Environment=RETENTION=14

# Resource hygiene — a backup must not starve the server.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=journal
StandardError=journal

# Light hardening. We still need access to the docker socket and the
# backup directory, so most filesystems can be read-only.
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/railways-secretariat /var/run/docker.sock /run/docker.sock
NoNewPrivileges=true
PrivateTmp=true
