# tool/secretariat.service
[Unit]
Description=Railway Secretariat API Server
After=network.target

[Service]
Type=simple
User=secretariat
Group=secretariat
WorkingDirectory=/opt/secretariat
ExecStart=/opt/secretariat/server_main
Environment=SECRETARIAT_SERVER_HOST=0.0.0.0
Environment=SECRETARIAT_SERVER_PORT=8080
Environment=SECRETARIAT_STORAGE_ROOT=/opt/secretariat/secretariat_data
Restart=always
RestartSec=5
StandardOutput=append:/opt/secretariat/logs/server.log
StandardError=append:/opt/secretariat/logs/error.log

[Install]
WantedBy=multi-user.target
